# 100-automount-flash.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sde1" SYMLINK+="usbflash", RUN+="/root/automount-flash.sh"
